# epg2xml.timer
[Unit]
Description=Run epg2xml daily

[Timer]
OnCalendar=daily
Persistent=true     
 
[Install]
WantedBy=timers.target
